# etc/rc.local
#! /bin/sh
#
# Stub for per machine changes to /etc/rc
#
# $Source: /d2/3.7/src/etc/RCS/rc.local,v $
# $Revision: 1.1 $
# $Date: 89/03/27 15:38:31 $
#
